# build/gitea/gitea.service.tmpl
[Unit]
Description=gitea - A painless self-hosted Git service
Documentation=https://gitea.io/
After=postgresql.target
Wants=network-online.target systemd-networkd-wait-online.service sshd.service

[Service]
Restart=on-failure
ExecStart={{ gitea_path }} web --port {{ gitea_port }} --config /home/{{ gitea_user }}/gitea/custom/conf/app.ini
User={{ gitea_user }}
Group={{ gitea_group }}
WorkingDirectory=/home/{{ gitea_user }}/gitea
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

