# roles/grafana/templates/jupyter.service.tmpl
[Unit]
Description=jupyter - Python notebook server
Documentation=http://jupyter.github.io/
After=network-online.target
Wants=systemd-networkd-wait-online.service

[Service]
Restart=on-abnormal
ExecStart=/home/{{ jupyter_user }}/.local/bin/jupyter notebook
User={{ jupyter_user }}
Group={{ jupyter_group }}
WorkingDirectory=/home/{{ jupyter_user }}
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

